# swap.service
[Unit]
Description=swap

[Service]
Type=forking
SmackProcessLabel=swap
ExecStart=/usr/bin/da_manager
