# swap.socket
[Unit]
Description=debug_socket

[Socket]
User=*
Group=*
ListenStream=/tmp/da.socket
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@

[Socket]
User=*
Group=*
ListenStream=/tmp/da_ui.socket
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@


[Install]
WantedBy=sockets.target
